# sched/launchd/radar.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.local.radar.weekly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string><string>python3</string>
    <string>cli/radar.py</string><string>run</string><string>--mode</string><string>weekly</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>WorkingDirectory</key><string>REPO_ABSOLUTE_PATH</string>
  <key>StandardOutPath</key><string>REPO_ABSOLUTE_PATH/.radar/logs/weekly.out.log</string>
  <key>StandardErrorPath</key><string>REPO_ABSOLUTE_PATH/.radar/logs/weekly.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OLLAMA_HOST</key><string>http://127.0.0.1:11434</string>
  </dict>
  <key>RunAtLoad</key><false/>
</dict>
</plist>